# backend/sastaspace-backend.service
[Unit]
Description=SastaSpace Backend API with DuckDB
After=network.target

[Service]
Type=simple
User=mkhare
WorkingDirectory=/home/mkhare/Development/sastaspace/backend
Environment="PATH=/home/mkhare/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DB_PATH=/home/mkhare/Development/sastaspace/backend/data/sastaspace.db"
ExecStart=/home/mkhare/.local/bin/uv run --directory /home/mkhare/Development/sastaspace/backend uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
